# scripts/retail-auto-redeploy.service
[Unit]
Description=Retail system Docker auto-redeploy watcher
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/workspaces/demand-forecasting-drift
ExecStart=/workspaces/demand-forecasting-drift/scripts/auto-redeploy.sh
Restart=always
RestartSec=15
Environment=REPO_DIR=/workspaces/demand-forecasting-drift
Environment=BRANCH=main
Environment=INTERVAL_SECONDS=60

[Install]
WantedBy=multi-user.target